# .config/systemd/user/set_time_wallpaper.service
[Unit]
Description=Set time wallpaper for sway

[Service]
Type=oneshot
ExecStart=/home/rizhiy/.bin/set_time_wallpaper_sway.sh
